# scraper/systemd/slot-scraper-bls.service
[Unit]
Description=Slot scraper (BLS International)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=scraper
WorkingDirectory=/opt/slot-scraper
EnvironmentFile=/opt/slot-scraper/.env
ExecStart=/usr/bin/node /opt/slot-scraper/scrape-bls.mjs
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/slot-scraper

[Install]
WantedBy=multi-user.target